# supervisor/start_vnc.conf
[program:start_vnc]
command=vncserver
autostart=true
autorestart=true
stderr_logfile=/var/log/vnc.err.log
stdout_logfile=/var/log/vnc.out.log